# systemd/hybrid_rasa_llm.service
[Unit]
Description=Hybrid Rasa LLM Flask Service
After=network.target

[Service]
ExecStart=/home/stjl0/livekit-travel-voice-assistant/llm_venv/bin/python /home/stjl0/livekit-travel-voice-assistant/hybrid_rasa_llm.py
WorkingDirectory=/home/stjl0/livekit-travel-voice-assistant
Restart=always
User=stjl0
Environment="PATH=/home/stjl0/livekit-travel-voice-assistant/llm_venv/bin"
StandardOutput=append:/home/stjl0/livekit-travel-voice-assistant/hybrid_rasa_llm.log
StandardError=append:/home/stjl0/livekit-travel-voice-assistant/hybrid_rasa_llm.log

[Install]
WantedBy=multi-user.target
