# livekit-travel-assistant/systemd/currency_converter.service
[Unit]
Description=Currency Converter Background Service
After=network.target

[Service]
User=stjl0
Group=stjl0
WorkingDirectory=/home/stjl0/Travel-Assistant-Application/livekit-travel-assistant
Environment="PATH=/home/stjl0/Travel-Assistant-Application/livekit-travel-assistant/llm_venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin"
ExecStart=/home/stjl0/Travel-Assistant-Application/livekit-travel-assistant/llm_venv/bin/python /home/stjl0/Travel-Assistant-Application/livekit-travel-assistant/currency_converter.py
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
